# webhook/webhook.service
[Unit]
Description=Alertmanager Webhook Receiver
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/webhook_server.py
Restart=always
User=root

[Install]
WantedBy=multi-user.target

